# claude-config/launchd/com.claude-learn.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-learn.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>claude --print "/learn --apply --cross-project --validate" >> "$HOME/Library/Logs/claude-learn/weekly.log" 2>&amp;1</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jasonjob/agents</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>5</integer>
        <key>Hour</key>
        <integer>22</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jasonjob/Library/Logs/claude-learn/weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jasonjob/Library/Logs/claude-learn/weekly.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
