# claude-config/launchd/com.claude-learn-gate-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-learn-gate-poller</string>
    <key>ProgramArguments</key>
    <!-- Logic lives in claude-config/scripts/learn-gate-poller.sh so it is
         readable/testable and editable without re-deploying this plist.
         `bash -lc` = login shell so git/gh credential helpers are on PATH. -->
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>exec "$HOME/agents/claude-config/scripts/learn-gate-poller.sh"</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jasonjob/agents</string>
    <key>StartInterval</key>
    <integer>43200</integer>
    <key>StandardOutPath</key>
    <string>/Users/jasonjob/Library/Logs/claude-learn/poller.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jasonjob/Library/Logs/claude-learn/poller.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
