# claude-config/launchd/com.claude-learn-gate-poller.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-learn-gate-poller</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>
          LOG="$HOME/Library/Logs/claude-learn/poller.log"
          mkdir -p "$(dirname "$LOG")"
          cd "$HOME/agents" \
          &amp;&amp; git pull --ff-only --quiet 2&gt;&gt;"$LOG" \
          &amp;&amp; python3 claude-config/scripts/telemetry_gate.py --verbose \
               &gt;&gt;"$LOG" 2&gt;&amp;1 \
          &amp;&amp; claude --print "/learn --apply --cross-project --validate" \
               &gt;&gt;"$LOG" 2&gt;&amp;1 \
          || echo "[$(date -Iseconds)] gate not tripped or pull failed — skip" \
               &gt;&gt;"$LOG"
        </string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jasonjob/agents</string>
    <key>StartInterval</key>
    <integer>43200</integer>
    <key>StandardOutPath</key>
    <string>/Users/jasonjob/Library/Logs/claude-learn/poller.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jasonjob/Library/Logs/claude-learn/poller.err</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
