# claude-config/launchd/com.coding-memory-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- coding-memory: daily personal ingest into the jns store. Logic lives in
     scripts/coding-memory-ingest.sh (editable without redeploying this plist).
     Default personal sources only (residency). NOT auto-installed by install.sh;
     deploy manually (see launchd/README.md). RunAtLoad false: first run is one
     interval after bootstrap. SSH key auth only — no credentials in this plist. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.coding-memory-ingest</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/jasonjob/agents/claude-config/scripts/coding-memory-ingest.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>86400</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/jasonjob/.claude/logs/coding-memory-ingest.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jasonjob/.claude/logs/coding-memory-ingest.log</string>
</dict>
</plist>
